# scripts/schedule/target-intel-weekly.timer
[Unit]
Description=Target Intelligence — 每周周报 timer（北京时间周一 07:00 = UTC 周日 23:00）
Requires=target-intel-weekly.service

[Timer]
OnCalendar=Sun *-*-* 23:00:00
Persistent=true
RandomizedDelaySec=600

[Install]
WantedBy=timers.target
